# src/pesign.service.in
[Unit]
Description=Pesign signing daemon

[Service]
PrivateTmp=true
PIDFile=@@RUNDIR@@/pesign.pid
ExecStart=/usr/bin/pesign --daemonize --nofork

[Install]
WantedBy=multi-user.target
